# biglinux-manual/etc/xdg/autostart/biglinux-manual.desktop
[Desktop Entry]
Name=biglinux-manual
Type=Application
Exec=/usr/bin/biglinux-manual-create-link.sh
NoDisplay=true
